# soccer_widget/.config/systemd/user/soccer_widget.timer
[Unit]
Description=Timer for soccer widget data update script

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
